# infrastructure/ikim_koesterlab/hosts/ikim-koesterlab-compute01/system_files/etc/profile.d/env.sh
if groups | grep koesterlab
then
    export TMPDIR=/mnt/tmp
fi
