# infrastructure/ikim_koesterlab/hosts/ikim-koesterlab-compute01/system_files/etc/profile.d/koesterlab_env.sh
if groups | grep koesterlab > /dev/null
then
    export TMPDIR=/mnt/tmp
    export SNAKEMAKE_PROFILE=default
    export SNAKEMAKE_OUTPUT_CACHE=/mnt/workspace/snakemake_cache
fi
