# infrastructure/koesterlab_ikim/hosts/ikim-koesterlab-compute01/system_files/etc/profile.d/env.sh
if getent group koesterlab | grep -qw "$USER"
then
    export TMPDIR=/mnt/tmp
fi
